# test/crontab/file1.cron.list
[32mellipsis.update[0m
  @reboot
  $HOME/.ellipsis/bin/ellipsis update >/dev/null 2>&1
[32mellipsis.test[0m
  @reboot
  echo "ellipsis.test"
[33mellipsis.disabled[0m
  @reboot
  echo "ellipsis.test"
