# test/crontab/file2.cron
# Ellipsis-Cron : ellipsis.test_new
@reboot echo "test"
