# test/crontab/file1.cron.clean
@reboot echo "test"


# Random comment!
